# remnav/metadata/cmd/gpsdol/gpsd.socket
[Socket]
# First blank ListStream clears the system defaults
ListenStream=
ListenStream=/var/run/gpsd.sock
ListenStream=0.0.0.0:2947
